# deploy/manager/systemd/caddy.service.d/realtime-me.conf
[Unit]
After=openvpn-server@realtime-me.service
Wants=openvpn-server@realtime-me.service

[Service]
EnvironmentFile=/etc/realtime-me/caddy.env
